# resources/dawn.conf
start on stopping update

script
  cd ~/daemon/app
  npm start
end script
