# resources/dawn.conf
start on stopping update
setuid vagrant

script
  cd /home/vagrant/daemon/app # FIXME
  npm start
end script
