# resources/hibike.conf
start on stopping update

script
  HOME_DIR=/home/vagrant # FIXME FIXME FIXME
  cd $HOME_DIR/hibike
  python hibike.py
end script
