# resources/hibike.conf
# we may not need this config file at all depending on how hibike is run
start on stopping update

script
  cd ~/hibike
  exec python hibike.py
end script
